# systemd/xibi-heartbeat.service
[Unit]
Description=Xibi Heartbeat Poller
After=network.target xibi-telegram.service
Wants=network.target

[Service]
Type=simple
User=%U
WorkingDirectory=%h/xibi
EnvironmentFile=%h/.xibi/secrets.env
ExecStart=%h/.local/bin/python -m xibi heartbeat \
    --config %h/.xibi/config.json \
    --workdir %h/.xibi
Restart=on-failure
RestartSec=30
TimeoutStopSec=120
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xibi-heartbeat

[Install]
WantedBy=default.target
