# systemd/xibi-telegram.service
[Unit]
Description=Xibi Telegram Adapter
After=network.target
Wants=network.target

[Service]
Type=simple
User=dlebron
WorkingDirectory=/home/dlebron/xibi
EnvironmentFile=/home/dlebron/.xibi/secrets.env
ExecStart=/home/dlebron/.local/bin/python -m xibi telegram \
    --config /home/dlebron/.xibi/config.json \
    --workdir /home/dlebron/.xibi
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xibi-telegram

[Install]
WantedBy=default.target
