# systemd/xibi-autoupdate.service
[Unit]
Description=Xibi Auto-Deploy (pull from GitHub and restart if changed)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/xibi
ExecStart=/bin/bash %h/xibi/scripts/xibi_deploy.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xibi-autoupdate
