# systemd/xibi-telegram.service
[Unit]
Description=Xibi Telegram Adapter
After=network.target
Wants=network.target

[Service]
Type=simple
WorkingDirectory=%h/xibi
EnvironmentFile=%h/.xibi/secrets.env
ExecStart=/usr/bin/python3 -m xibi --config %h/.xibi/config.json --workdir %h/.xibi telegram
Restart=on-failure
RestartSec=10
TimeoutStopSec=300
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xibi-telegram

[Install]
WantedBy=default.target
