# scripts/xibi-deploy.service
[Unit]
Description=Xibi auto-deploy (poll origin/main, pull, restart)

[Service]
Type=oneshot
ExecStart=%h/xibi/scripts/deploy.sh
Environment=HOME=%h
